# dist/form2mail.service
[Unit]
Description=formto.email
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

[Service]
Restart=on-failure
StartLimitInterval=360
StartLimitBurst=5

; User and group the process will run as.
User=caddy
Group=caddy
WorkingDirectory=/var/www/form2mail

ExecStart=/var/www/formto.email/form2mail
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=multi-user.target
